# centos/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D

[program:wait]
command=/bin/bash -c "tail -f /dev/null"
